# nvidia-core-temp.service
[Unit]
Description=Update NVIDIA GPU Core Temperature (NVML → hwmon)

[Service]
Type=oneshot
ExecStart=/usr/libexec/nvidia-core-temp-helper
